# weather-watcher.service
[Unit]
Description=Weather Watcher - AccuWeather Radar Display
After=network.target graphical.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/weather-watcher
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target

